# deploy/postgres.container
# Postgres database for bits
# Security: Internal network only, no external port exposure

[Unit]
Description=Bits Postgres database
After=network-online.target

[Container]
Image=docker.io/library/postgres:17-alpine
ContainerName=bits-postgres

# Network: internal only, no published ports
Network=bits.network

# Persistent data
Volume=bits-pgdata.volume:/var/lib/postgresql/data:Z

# Database credentials for initial bootstrap (ragenix)
# Required: POSTGRES_USER, POSTGRES_PASSWORD, POSTGRES_DB
EnvironmentFile=/run/agenix/bits-postgres-environment-file

# Security hardening
ReadOnly=true
NoNewPrivileges=true
DropCapability=ALL
AddCapability=CAP_CHOWN CAP_FOWNER CAP_SETGID CAP_SETUID CAP_DAC_OVERRIDE

# Postgres needs to write to these tmpfs locations
Tmpfs=/tmp:size=64m,mode=1777
Tmpfs=/run:size=64m,mode=755
Tmpfs=/var/run/postgresql:size=64m,mode=755,uid=999,gid=999

# Resource limits
Memory=512M
MemorySwap=512M

# Health check
HealthCmd=pg_isready -U bits -d bits
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=default.target
